# scripts/systemd/a2a2h-cache-keepalive.service
[Unit]
Description=A2A2H prompt-cache keep-alive (pings live sessions to avoid cold re-bootstrap)

[Service]
Type=oneshot
ExecStart=/opt/a2a2h/scripts/cache-keepalive.sh
StandardOutput=journal
StandardError=journal
