# scripts/systemd/a2a2h-cache-keepalive.timer
[Unit]
Description=Fire cto-cache-keepalive every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=cto-cache-keepalive.service

[Install]
WantedBy=timers.target
