# scripts/systemd/a2a2h-cache-keepalive.service
[Unit]
Description=CTO prompt-cache keep-alive (pings live sessions to avoid cold re-bootstrap)

[Service]
Type=oneshot
ExecStart=/opt/cto/scripts/cache-keepalive.sh
StandardOutput=journal
StandardError=journal
